# conf/supervisor/xunsearch.ini
[program:xunsearch]
command=xs-ctl -b inet start