# conf/supervisor/xunsearch.ini
[program:xunsearch]
command=xs-ctl -b inet start
autostart=true
stderr_logfile=/opt/websrv/logs/xunsearch_stderr.log
stdout_logfile=/opt/websrv/logs/xunsearch_stdout.log